# crontab
0 0 * * * sh daily.sh
0 * * * * sh hourly.sh
